# lib/systemd/system/mpd-sync.service
[Unit]
Description=mpd disk sync
After=mpd.service
[Service]
ExecStart=/usr/local/bin/mpd-sync

[Install]
WantedBy=multi-user.target
