# lib/systemd/system/nohdmi.service
[Unit]
Description=Disable HDMI Output

[Service]
ExecStart=/usr/bin/tvservice -o
Type=oneshot

[Install]
WantedBy=multi-user.target
